# orca_vision/scripts/stereo.service
[Unit]
Description=Stereo camera

[Service]
User=root
ExecStart=/home/pi/stereo.sh

[Install]
WantedBy=multi-user.target
